# safe-shutdown.service
[Unit]
Description=RetroFlag Pi Case safe shutdown

[Service]
ExecStart=/opt/retroflag-picase/service.py

[Install]
WantedBy=multi-user.target
